# packaging/systemd/jvm-sentinel-tomcat.service
[Unit]
Description=JVM Sentinel
After=network.target

[Service]
Type=simple
User=tomcat
Group=tomcat
ExecStart=/usr/local/bin/jvm-sentinel -config /etc/jvm-sentinel/tomcat.properties
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

# Keep Tomcat paths explicit for catalina.sh and custom scripts.
Environment=JAVA_HOME=/usr/lib/jvm/default-java
Environment=CATALINA_HOME=/opt/tomcat
Environment=CATALINA_BASE=/opt/tomcat
Environment=CATALINA_PID=/opt/tomcat/temp/tomcat.pid

# Keep the Go sentinel small. These are ignored by the JVM.
Environment=GOMEMLIMIT=24MiB
Environment=GOGC=50

# Basic hardening. Adjust ReadWritePaths for your Tomcat and backup paths.
# Do not add MemoryMax here when catalina.sh starts Tomcat directly; Java child
# processes may stay in this service cgroup and inherit the cap.
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ReadWritePaths=/opt/tomcat /var/backups/jvm-sentinel /var/log/jvm-sentinel

[Install]
WantedBy=multi-user.target
